# roles/supervisord/templates/supervisord.ini
[program:{{ project_name }}]
directory={{ project_dir }}/{{ project_name }}
command={{ java_dir }} {{ java_option_before }} -jar -Dserver.port="{{ app_port }}" {{ jvm_option }} {{ project_dir }}/{{ project_name }}/{{ project_name }}.jar --logging.path={{ logs_dir }}/{{ project_name }}/ {{ java_option_after }} {{ project_name_config }}
user={{ app_user }}
autostart=true
autorestart=unexpected
startsecs=15
exitcodes=2
startretries=3
environment=